# deploy/eda-bot.service
# Заготовка systemd-юнита для Fornex. Пути уточнить при деплое.
# /etc/systemd/system/eda-bot.service

[Unit]
Description=Eda cooking bot (@edaforme_bot)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/eda
ExecStart=/opt/eda/.venv/bin/python -m bot.main
Restart=on-failure
RestartSec=5
# .env лежит в /opt/eda/.env — читается самим ботом, EnvironmentFile не нужен

[Install]
WantedBy=multi-user.target
